# openrouter-monitor.service
[Unit]
Description=OpenRouter Discord Monitor
Documentation=https://github.com/Yokonad/orpdisc
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/yokonad/yndhome/yndthings/orpdic
EnvironmentFile=/etc/openrouter-monitor.env
ExecStart=/home/yokonad/yndhome/yndthings/orpdic/monitor
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openrouter-monitor

[Install]
WantedBy=multi-user.target